# crontab
22 3    * * *   www-data  php /var/www/html/admin.php clean_expired
22 3    * * *   www-data  php /var/www/html/admin.php clean_async

